# docker/supervisord.conf
[supervisord]
nodaemon=true
user=appuser
logfile=/data/logs/supervisord.log
pidfile=/data/logs/supervisord.pid
childlogdir=/data/logs


[program:hermes-gateway]
command=python -m uvicorn app:gateway_app --host 0.0.0.0 --port 8001
directory=/app
autostart=true
autorestart=true
startretries=3
redirect_stderr=false
stdout_logfile=/data/logs/hermes-gateway.log
stderr_logfile=/data/logs/hermes-gateway.err.log
environment=PYTHONPATH="/app"

[program:fastapi]
command=python -m uvicorn app:fastapi_app --host 0.0.0.0 --port 8000
directory=/app
autostart=true
autorestart=true
startretries=3
redirect_stderr=true
stdout_logfile=/data/logs/fastapi.log
stderr_logfile=/data/logs/fastapi.err.log

[program:streamlit]
command=streamlit run app.py --server.port 8501 --server.address 0.0.0.0 --server.headless true
directory=/app
autostart=true
autorestart=true
startretries=3
redirect_stderr=true
stdout_logfile=/data/logs/streamlit.log
stderr_logfile=/data/logs/streamlit.err.log